# panopticon_capture.timer
[Unit]
Description=Panopticon image capture

[Timer]
# Run on the minute, every minute
OnCalendar=*-*-* *:*:01

[Install]
WantedBy=timers.target
